# systemd/podkast.radiorevolt.no-images.timer
[Unit]
Description=Run image processing twice every hour

[Timer]
OnCalendar=*:25,59:00
Persistent=true

[Install]
WantedBy=timers.target
